# etc/cron.d/check-network
#Network Check
*/5 * * * * root /usr/local/bin/check-network >/dev/null 2>&1
