# etc/cron.d/yum-update
0 2 * * * root yum update -y > /dev/null 2>&1
