# etc/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=41443
